# redshift-gtk.service
[Service]
ExecStart=/usr/bin/redshift-gtk
