# redshift-gtk.service.d/display.conf
[Service]
# For some reason, this only seems to fix the crash problem if it's in a
# subdirectory.
Environment=DISPLAY=:0
